# 32/setup/files/supervisor-files/vnc-supervisord.conf
[supervisord]

[program:x11vnc]
command=x11vnc -ncache 10 -ncache_cr -xkb -forever -repeat -capslock -usepw -create
